# crontab
SHELL=/bin/bash
* * * * * root sh /run.sh 2>&1

